# viz.service
[Unit]
Description=Viz

[Service]
User=pi
Environment=SDL_VIDEODRIVER=kmsdrm
ExecStart=/home/pi/viz
Restart=no

[Install]
WantedBy=multi-user.target
